# roles/restic/templates/alert.service.j2
{{ ansible_managed | comment }}
[Unit]
Description=Restic alert
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/restic-alert.sh

[Install]
WantedBy=multi-user.target
