# roles/restic/templates/forget.timer.j2
{{ ansible_managed | comment }}
[Unit]
Description=Restic forget

[Timer]
OnCalendar={{ restic_forget_timer }}

[Install]
WantedBy=timers.target
